# after_install/jarvis.service
[Unit]
Description=Jarvis AI Assistant
After=network.target sound.target

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/jarvis-assistant
Environment="PATH=/home/%i/jarvis-assistant/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/%i/jarvis-assistant/venv/bin/python /home/%i/jarvis-assistant/jarvis.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
